# etc/systemd/system/mycroft.service
[Unit]
Description=MyCroft
After=multi-user.target
Requires=sound.target
After=sound.target

[Service]
User=pradu
Type=forking
ExecStart=/home/pradu/mycroft/start-mycroft.sh all
ExecStop=/home/pradu/mycroft/stop-mycroft.sh all

[Install]
WantedBy=graphical.target
Wants=pulseaudio.service devsetup.service
